# crontab
TZ=Asia/Seoul
PATH=/usr/local/bin
0 6 * * * python -u /app/main.py >> /var/log/cron.log 2>&1
