# crontab
CRON_TZ=Asia/Seoul
0 6 * * * python -u /app/main.py >> /var/log/cron.log
